# vm-rollback.timer
[Unit]
Description=Run Proxmox VM rollback every 3 days

[Timer]
OnBootSec=10m
OnUnitActiveSec=3d
Persistent=true

[Install]
WantedBy=timers.target
